# use/ansible/roles/js/templates/etc/profile.d/phantomjs.sh
#!/usr/bin/env bash
# {{ ansible_managed }}

export PATH=$PATH:"{{ workstation_software_root_directory }}/phantomjs/bin"