# use/ansible/roles/java/templates/etc/profile.d/ant.sh
#!/usr/bin/env bash
# {{ ansible_managed }}

export ANT_HOME={{ workstation_software_root_directory }}/ant
export PATH=$PATH:$ANT_HOME/bin
